# net-misc/smcroute/files/smcroute.initd
#!/sbin/openrc-run
# Copyright 1999-2026 Gentoo Authors
# Distributed under the terms of the GNU General Public License v2
# shellcheck shell=sh disable=SC2034

command="/usr/sbin/smcrouted"
pidfile="/var/run/smcroute.pid"
command_args="-p smcroute:smcroute -P ${pidfile} -N"

depend() {
	need localmount net
	use logger
}
