# net-nds/authelia/files/authelia.initd
#!/sbin/openrc-run

command=/usr/bin/authelia
command_args="${authelia_opts:---config /etc/authelia/configuration.yml}"
command_background=yes
command_user=authelia:authelia
make_pidfile=yes
pidfile=/run/authelia.pid
